# deploy/upwork-dashboard.service
[Unit]
Description=Upwork Scraper Streamlit Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=npc
Group=npc
WorkingDirectory=/home/npc/upwork-scrap

# Load DATABASE_URL for cloud DB access
EnvironmentFile=/home/npc/upwork-scrap/.env

ExecStart=/home/npc/upwork-scrap/.venv/bin/streamlit run dashboard/app.py \
    --server.port=8501 \
    --server.address=0.0.0.0 \
    --server.headless=true

Restart=always
RestartSec=10

StandardOutput=append:/var/log/upwork-scrap/dashboard.log
StandardError=append:/var/log/upwork-scrap/dashboard.log

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/home/npc/upwork-scrap/data /var/log/upwork-scrap

[Install]
WantedBy=multi-user.target
